# deploy/systemd/orchestrai-validator.service
[Unit]
Description=OrchestrAI - Validator Engine (Data Ops Validator)
After=network.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=all
User=orchestrai
Group=orchestrai
WorkingDirectory=/var/www/orchestrai
EnvironmentFile=/var/www/orchestrai/.env
ExecStart=/var/www/orchestrai/venv/bin/gunicorn \
    --config /var/www/orchestrai/gunicorn.conf.py \
    --bind 127.0.0.1:8000 \
    app.main:app
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orchestrai-validator

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/orchestrai/data /var/www/orchestrai/validated_runs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
